# cockroach.conf
[program:cockroachdb]
directory=/mnt/vol1
command=/home/rock64/spot/bin/huhuladb.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/huhuladb.err.log
stdout_logfile=/var/log/huhuladb.out.log

